# templates/supervisord.conf
[supervisord]
; just log to stdout so you can see everything in the fly monitoring
logfile=/dev/stdout 
logfile_maxbytes=0  
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:agent-api]
command=AGENT_RUN_COMMAND
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xmtp-service]
command=node xmtp-service.js
directory=/xmtp-service
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0